# deploy/13flow-stats.service
[Unit]
Description=Generate the private 13FLOW GoAccess report
After=apache2.service
ConditionPathIsDirectory=/var/www/html/13flow-stats
ConditionPathExists=/var/log/apache2/13flow_access.log

[Service]
Type=oneshot
User=flowstats
Group=flowstats
SupplementaryGroups=adm
WorkingDirectory=/var/lib/13flow-stats
Environment=HOME=/var/lib/13flow-stats
Environment=LC_TIME=C
UMask=0077
ExecStart=/usr/local/libexec/13flow-generate-stats

NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
PrivateNetwork=yes
PrivateIPC=yes
ProtectSystem=strict
ProtectHome=yes
ProtectProc=invisible
ProcSubset=pid
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
KeyringMode=private
RemoveIPC=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX
IPAddressDeny=any
StateDirectory=13flow-stats
StateDirectoryMode=0700
ReadOnlyPaths=/var/log/apache2
ReadWritePaths=/var/www/html/13flow-stats /var/lib/13flow-stats
InaccessiblePaths=-/etc/13flow -/etc/apache2/13flow-stats.htpasswd -/var/lib/13flow -/var/lib/13flow-pro -/var/lib/13flow-mcp
MemoryMax=256M
TasksMax=32
Nice=10
IOSchedulingClass=idle
