# deploy/13flow-pro-workspace-snapshot.timer
[Unit]
Description=Schedule 13FLOW Pro workspace signal snapshots

[Timer]
OnCalendar=*-*-* 07:45:00
Persistent=true
RandomizedDelaySec=1200

[Install]
WantedBy=timers.target
